# scripts/launchd/org.crghouston.opportunity-scan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Weekly Opportunity Scan — runs functions/_lib/scan-pipeline.js via
  scripts/run-scan.mjs every Monday at 7:00 AM LOCAL time. launchd follows the
  machine's timezone and DST on its own, so no UTC drift to reason about.

  Install / update / remove: see scripts/launchd/README.md.

  Notes baked into the choices below:
  - Absolute node path — launchd gives a job a minimal PATH, so `node` alone
    would not resolve. Verify the path with `which node` and update if it moves.
  - RunAtLoad is false — we do NOT want a scan every time the agent is loaded
    (login / reboot); only on the Monday 2pm schedule.
  - If the Mac is asleep at 2pm but awake later, launchd runs the missed job on
    wake. If it was powered off or sitting at the FileVault/login screen, the
    run is simply skipped — an accepted loss (the feed is nice-to-have).
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>org.crghouston.opportunity-scan</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/Users/operacha/crg2026-app/scripts/run-scan.mjs</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/operacha/crg2026-app</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string>/Users/operacha/crg2026-app/logs/scan.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/operacha/crg2026-app/logs/scan.err.log</string>
</dict>
</plist>
